# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1

# 每10分钟同步个人日记
*/10 9-22 * * *  root bash /var/www/my-notes/.script/sync.sh >> /dev/null 2>&1
